# deploy/systemd/boost-match-backend.service
[Unit]
Description=BoostGateway v2 Matchmaking Backend
Documentation=https://github.com/HoneyBury/BoostAsioDemo
After=network.target

[Service]
Type=simple
User=boost-gateway
Group=boost-gateway
Environment=MATCH_PORT=9304
ExecStart=/usr/local/bin/v2_match_backend
WorkingDirectory=/var/lib/boost-gateway/match
Restart=on-failure
RestartSec=2s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boost-match-backend

RuntimeDirectory=boost-gateway/match
RuntimeDirectoryMode=0750
StateDirectory=boost-gateway/match
StateDirectoryMode=0750

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/boost-gateway/match /var/log/boost-gateway

[Install]
WantedBy=multi-user.target
