# deploy/systemd/boost-login-backend.service
[Unit]
Description=BoostGateway v2 Login Backend
Documentation=https://github.com/HoneyBury/BoostAsioDemo
After=network.target

[Service]
Type=simple
User=boost-gateway
Group=boost-gateway
ExecStart=/usr/local/bin/v2_login_backend 9202
WorkingDirectory=/var/lib/boost-gateway/login
Restart=on-failure
RestartSec=2s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boost-login-backend

RuntimeDirectory=boost-gateway/login
RuntimeDirectoryMode=0750
StateDirectory=boost-gateway/login
StateDirectoryMode=0750

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/boost-gateway/login /var/log/boost-gateway

[Install]
WantedBy=multi-user.target
