# deploy/systemd/boost-leaderboard-backend.service
[Unit]
Description=BoostGateway v2 Leaderboard Backend
Documentation=https://github.com/HoneyBury/BoostAsioDemo
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=boost-gateway
Group=boost-gateway
Environment=LEADERBOARD_PORT=9305
Environment=REDIS_HOST=127.0.0.1
Environment=REDIS_PORT=6379
ExecStart=/usr/local/bin/v2_leaderboard_backend
WorkingDirectory=/var/lib/boost-gateway/leaderboard
Restart=on-failure
RestartSec=2s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boost-leaderboard-backend

RuntimeDirectory=boost-gateway/leaderboard
RuntimeDirectoryMode=0750
StateDirectory=boost-gateway/leaderboard
StateDirectoryMode=0750

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/boost-gateway/leaderboard /var/log/boost-gateway

[Install]
WantedBy=multi-user.target
